# resources/supervisord.conf
[supervisord]
nodaemon=true

[program:xlr]
command=/opt/xlr/server/bin/run.sh
