# resources/xlr.conf
[program:xlr]
command=/opt/xlr/server/bin/run.sh

[program:log_redirect]
command=sh -c "sleep 10;tail -f /opt/xlr/server/log/xl-release.log"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
